# scripts/zapret-netwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>zapret-netwatch</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/zapret/local-tools/zapret-netwatch.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>15</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Library/Application Support/Zapret/logs/netwatch-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/Application Support/Zapret/logs/netwatch-stderr.log</string>
</dict>
</plist>
